# sherssr.conf
[program:sherssr]
directory=/var/www/sher.gtfb.tech
command=node server-dist/index.js
numprocs=1
autostart=true
autorestart=true
user=www-data